# com.lifeos.bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!--
  This is a TEMPLATE. Run ./configure.sh to generate a personalized copy
  installed at ~/Library/LaunchAgents/com.lifeos.bot.plist
-->
<dict>
  <key>Label</key>
  <string>com.lifeos.bot</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>__LIFEOS_DIR__/run_lifeos_bot.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__LIFEOS_DIR__</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>__HOME__/.lifeos/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>__HOME__/.lifeos/launchd.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
</dict>
</plist>
